# etc/runit/sv/NetworkManager-dispatcher/run
#!/bin/sh
sv check NetworkManager >/dev/null || exit 1
exec /usr/lib/nm-dispatcher > /dev/null 2>&1
