# etc/runit/sv/intel_hwp/run
#!/bin/sh
exec x86_energy_perf_policy --hwp-enable balance-power
