# dotfiles/.config/runit/sv/psd/finish
#!/bin/sh

exec profile-sync-daemon unsync
