# dotfiles/.config/runit/sv/psd/run
#!/bin/sh

profile-sync-daemon sync
exec chpst -b profile-sync-daemon pause
