# dotfiles/.config/runit/sv/psd/run
#!/bin/sh

profile-sync-daemon startup
exec chpst -b psd pause
